# alpine-supervisord/ping.conf
[program:ping]
command=ping tut.by
user=root
